# dot_config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland compositor session
Documentation=man:systemd.special(7)
Wants=graphical-session.target
BindsTo=graphical-session.target
After=graphical-session.target
