# dot_config/systemd/user/solaar.service
[Unit]
Description=Solaar

[Service]
ExecStart=/usr/bin/solaar

[Install]
WantedBy=default.target


